# deploy/zeroclaw.service
[Unit]
Description=ZeroClaw Autonomous Agent Runtime
Documentation=https://github.com/zeroclaw-labs/zeroclaw
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/bac/src/zeroclaw
ExecStart=/home/%i/bac/src/zeroclaw/target/release/zeroclaw daemon
Restart=on-failure
RestartSec=5

# Environment variables
Environment=ZEROCLAW_PROVIDER=openrouter
Environment=TELEGRAM_BOT_TOKEN=
Environment=TELEGRAM_DEFAULT_CHAT_ID=
Environment=WHATSAPP_ACCESS_TOKEN=
Environment=WHATSAPP_PHONE_NUMBER_ID=
Environment=WHATSAPP_DEFAULT_TO=

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/%i/.zeroclaw /home/%i/bac/.tmp/zeroclaw

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
